# supervisord.conf
[supervisord]

[program:reminder.citadel]
command=reminder.citadel.py
autorestart=true