# supervisor.conf
[program:reminder.citadel]
command=/opt/reminder.citadel/python3-venv/bin/python3 /opt/reminder.citadel/main.py
autorestart=true